# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log

[program:uvicorn]
command=/usr/local/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --log-level debug
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/uvicorn.err.log
stdout_logfile=/var/log/supervisor/uvicorn.out.log
stderr_logfile_maxbytes=50MB
stdout_logfile_maxbytes=50MB
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/backend",DATABASE_URL="sqlite:///./pim.db",SECRET_KEY="%(ENV_SECRET_KEY)s",ADMIN_EMAIL="%(ENV_ADMIN_EMAIL)s",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s",CORS_ORIGINS="%(ENV_CORS_ORIGINS)s",ACCESS_TOKEN_EXPIRE_MINUTES="%(ENV_ACCESS_TOKEN_EXPIRE_MINUTES)s",REFRESH_TOKEN_EXPIRE_MINUTES="%(ENV_REFRESH_TOKEN_EXPIRE_MINUTES)s"
